# server.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn2.sock

[Install]
WantedBy=sockets.target
